# src/solstone_linux/solstone-linux.service.in
[Unit]
Description=Solstone Linux Desktop Observer
After=graphical-session.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart={BINARY} run
PassEnvironment=DISPLAY WAYLAND_DISPLAY DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR XDG_CURRENT_DESKTOP XAUTHORITY XDG_SESSION_TYPE
Environment=PATH={PATH}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
